# server/node.conf
[program:node]
command=node /AMHub/server/server.js
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/node.err.log
stdout_logfile=/var/log/supervisor/node.out.log